# ollama.service
[Unit]
Description=Ollama LLM Server
After=network.target

[Service]
Type=simple
User=marcin
ExecStart=/usr/local/bin/ollama serve
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
